# deploy/systemd/ryazanvpn-node-agent.service
[Unit]
Description=RyazanVPN Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ryazanvpn
Group=ryazanvpn
WorkingDirectory=/opt/ryazanvpn
EnvironmentFile=/etc/ryazanvpn/node-agent.env
ExecStart=/opt/ryazanvpn/bin/node-agent
Restart=always
RestartSec=2
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
